# target/linux/sg200x/image/tmp/root/usr/lib/opkg/info/base-files.list
/mnt/system/lib/libcvi_vqe.so
/mnt/system/usr/bin/sample_integ
/mnt/system/ko/3rd/ts_gslX680.ko
/etc/openwrt_version
/etc/banner.failsafe
/mnt/system/usr/bin/sample_osdc
/mnt/data/sensor_cfg_OV5647_J2.ini
/etc/resolv.conf
/etc/board.d/01_leds
/mnt/system/lib/libini.so
/etc/init.d/done
/bin/config_generate
/var
/mnt/system/ko/cv181x_rtos_cmdqu.ko
/mnt/system/usb-host.sh
/mnt/system/usb-rndis.sh
/lib/preinit/00_preinit.conf
/mnt/system/ko/loadsystemko.sh
/etc/init.d/system
/mnt/system/usr/bin/sample_or
/etc/preinit
/lib/preinit/79_move_config
/mnt/system/firmware/aic8800/fw_adid_8800d80_u02.bin
/mnt/system/usr/bin/sample_dma
/mnt/system/lib/libz.so
/mnt/cfg/param/cvi_sdr_bin_OV5647.bin
/mnt/system/lib/libaaccomm2.so
/mnt/system/usr/bin/sample_magandang
/etc/localtime
/mnt/system/usr/bin/sample_cannyhysedge
/lib/upgrade/platform.sh
/mnt/system/lib/libaacenc2.so
/mnt/system/usr/bin/sample_venc
/mnt/system/usr/bin/sample_vio
/etc/dnsmasq.conf
/etc/shells
/mnt/system/usr/bin/sample_gmm2
/lib/functions/caldata.sh
/lib64
/lib/preinit/70_initramfs_test
/etc/rc.button/reboot
/mnt/system/lib/libsns_sc200ai.so
/mnt/system/usb.sh
/etc/device_info
/mnt/system/usr/bin/sample_region
/usr/lib64
/mnt/data/sensor_cfg_SC200AI.ini
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/mnt/system/ko/cv181x_jpeg.ko
/mnt/system/usr/bin/sample_query
/mnt/system/lib/libcvi_RES1.so
/lib/preinit/10_indicate_preinit
/mnt/system/lib/libvenc.so
/usr/libexec/login.sh
/mnt/system/ko/lzo_compress.ko
/usr/bin/bt_test
/mnt/system/ko/cv181x_dwa.ko
/mnt/system/firmware/aic8800/aic_userconfig_8800d80.txt
/mnt/cfg/param/cvi_sdr_bin_GC2083
/mnt/system/ko/cv181x_rtc.ko
/mnt/system/usr/bin/sample_sobel
/mnt/system/usr/bin/sample_hist
/mnt/system/usr/bin/sample_csc
/mnt/system/ko/cv181x_pwm.ko
/mnt/system/lib/libsns_sc035hgs.so
/mnt/system/usr/bin/sample_overlay
/lib/libatomic.so.1
/mnt/system/ko/cv181x_ive.ko
/mnt/system/ko/lzo.ko
/mnt/system/lib/libawb.so
/mnt/system/usr/bin/sample_lbp
/mnt/system/firmware/aic8800/fmacfwbt_8800d80_u02.bin
/etc/fstab
/usr/lib/os-release
/etc/sysctl.conf
/mnt/system/usr/bin/sample_map
/sbin/sysupgrade
/mnt/system/usr/bin/sample_dsi
/mnt/system/usr/bin/sample_gmm
/root/network-init.sh
/lib/preinit/99_10_run_init
/mnt/system/usr/bin/sample_vdecvo
/mnt/system/lib/libisp_algo.so
/mnt/system/usr/bin/sample_sub
/etc/run_usb.sh
/mnt/system/ko/cv181x_vpss.ko
/sbin/led.sh
/mnt/system/ko/zram.ko
/etc/iproute2/rt_protos
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/mnt/system/duo-init.sh
/mnt/system/lib/libtinyalsa.so
/mnt/system/usr/bin/sample_cvg
/mnt/system/ko/cvi_vc_driver.ko
/mnt/system/ko/lzo-rle.ko
/mnt/system/lib/libisp.so
/mnt/system/ko/cv181x_vi.ko
/etc/shadow
/mnt/system/usr/bin/sample_stcandicorner
/etc/rc.common
/mnt/system/lib/libcnpy.so
/etc/uci-defaults/10_migrate-shadow
/lib/preinit/02_sysinfo
/etc/board.d/99-default_network
/mnt/system/ko/3rd/cvi_wiegand_gpio.ko
/mnt/system/lib/libaacsbrdec2.so
/lib/preinit/10_indicate_failsafe
/mnt/system/firmware/aic8800/lmacfw_rf_8800d80_u02.bin
/mnt/system/lib/libraw_replay.so
/etc/os-release
/etc/board.d/02_network
/etc/network-init.sh
/mnt/system/lib/libz.so.1.2.11
/mnt/system/usr/bin/sample_thresh_S16
/mnt/system/usr/bin/sample_cannyedge
/mnt/system/usr/bin/sample_ncc
/mnt/system/ko/efivarfs.ko
/mnt/system/lib/libae.so
/mnt/system/ko/cvi_fb.ko
/lib/preinit/80_mount_root
/mnt/system/lib/libsns_ov5647.so
/etc/init.d/userinit
/mnt/system/lib/libopencv_core.so.3.2
/lib/libstdc++.so.6
/mnt/system/lib/libz.so.1
/etc/rc.local
/lib/functions/uci-defaults.sh
/mnt/system/usr/bin/sample_xor
/etc/init.d/led
/etc/ethers
/mnt/system/lib/libvpu.so
/mnt/system/ko/3rd/gt9xx.ko
/usr/bin/wifi_test
/etc/init.d/gpio_switch
/lib/functions/network.sh
/etc/init.d/sysctl
/mnt/system/lib/libsns_full.so
/mnt/system/lib/libsns_gc2053.so
/sbin/hotplug-call
/usr/libexec/validate_firmware_image
/mnt/system/lib/libsample.so
/mnt/system/lib/libaacdec2.so
/mnt/system/usr/bin/sample_16bitto8bit
/mnt/system/lib/libaf.so
/mnt/system/ko/cv181x_sys.ko
/mnt/system/usr/bin/sample_thresh_U16
/mnt/system/usr/bin/sample_erode
/etc/rc.button/failsafe
/mnt/system/lib/libmipi_tx.so
/etc/shinit
/mnt/system/ko/cv181x_saradc.ko
/lib/upgrade/fwtool.sh
/mnt/data/sensor_cfg_SC035HGS.ini
/etc/services
/mnt/system/usr/bin/sample_dilate
/mnt/cvimodel/scrfd_768_432_int8_1x.cvimodel
/lib/functions/system.sh
/mnt/system/lib/libcvi_ispd2.so
/mnt/system/lib/libcvi_ive.so
/mnt/system/lib/libcvi_tdl_app.so
/mnt/system/lib/libcviruntime.so
/bin/board_detect
/mnt/system/lib/libcvi_rtsp.so
/lib/upgrade/common.sh
/usr/bin/sg200x_pinmux
/etc/opkg/distfeeds.conf
/mnt/system/ko/cv181x_wdt.ko
/mnt/system/ko/aic8800_btlpm.ko
/mnt/system/ko/aic8800_bsp.ko
/etc/rc.button/power
/etc/sysupgrade.conf
/mnt/system/usr/bin/sample_filter
/mnt/system/usr/bin/sample_add
/mnt/system/usr/bin/sample_fisheye
/mnt/system/ko/cv181x_rgn.ko
/mnt/system/usr/bin/sample_resize
/mnt/system/usr/bin/sample_vcodec
/mnt/system/lib/libopencv_imgproc.so.3.2
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/mnt/system/lib/libraw_dump.so
/mnt/system/ko/cviaudio_core.ko
/etc/uhubon.sh
/lib/preinit/40_run_failsafe_hook
/lib/upgrade/do_stage2
/etc/group
/etc/rc.button/reset
/mnt/system/usr/bin/sample_sad
/mnt/system/usr/bin/sample_framediffmotion
/usr/lib/libwiringx.so
/lib/preinit/30_failsafe_wait
/lib/upgrade/stage2
/etc/opkg/keys/72d7d0d80f676aac
/mnt/data/sensor_cfg_GC2083.ini
/sbin/firstboot
/mnt/system/lib/libcvi_bin.so
/usr/lib64xthead/lp64d/libc.so
/mnt/system/lib/libcvi_tdl.so
/etc/uci-defaults/13_fix-group-user
/mnt/system/usb-ncm.sh
/mnt/system/lib/libcvi_audio.so
/etc/hotplug.d/net/00-sysctl
/etc/hosts
/mnt/system/ko/cv181x_clock_cooling.ko
/mnt/system/ko/cv181x_tpu.ko
/lib/functions.sh
/mnt/system/blink.sh
/mnt/system/lib/libmisc.so
/mnt/system/usr/bin/sample_thresh
/etc/iproute2/rt_tables
/mnt/system/firmware/aic8800/fw_patch_table_8800d80_u02.bin
/mnt/system/usr/bin/sensor_test
/mnt/system/usr/bin/sample_bgmodel
/lib/functions/leds.sh
/mnt/system/usr/bin/sample_cipher
/mnt/system/ko/cv181x_vcodec.ko
/mnt/system/usr/bin/sample_vdec
/rom/note
/etc/inittab
/mnt/system/usr/bin/sample_audio
/etc/openwrt_release
/mnt/system/lib/libosdc.so
/sbin/wifi
/mnt/system/lib/libsns_gc4653.so
/mnt/system/lib/libcvi_ssp.so
/mnt/system/lib/libopencv_imgcodecs.so.3.2
/etc/profile
/mnt/system/lib/libcvikernel.so
/mnt/system/usr/bin/camera-test.sh
/mnt/system/lib/libcvi_bin_isp.so
/mnt/system/ko/snsr_i2c.ko
/mnt/system/firmware/aic8800/fw_patch_8800d80_u02.bin
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/mnt/system/ko/cv181x_mon.ko
/mnt/system/lib/libcli.so
/etc/TZ
/mnt/system/usr/bin/sample_normgrad
/mnt/system/ko/cv181x_mipi_tx.ko
/etc/init.d/umount
/lib/functions/service.sh
/mnt/system/ko/cv181x_base.ko
/mnt/system/lib/libsns_gc2083.so
/mnt/system/usr/bin/sample_vi_fd
/mnt/system/usr/bin/sample_gradfg
/lib/upgrade/keep.d/base-files-essential
/lib/functions/migrations.sh
/etc/sysctl.d/10-default.conf
/mnt/system/firmware/aic8800/fmacfw_8800d80_u02.bin
/etc/passwd
/lib/preinit/02_default_set_state
/mnt/system/lib/libvdec.so
/etc/iproute2/ematch_map
/mnt/system/usr/bin/sample_test_tp2863_tp2803
/mnt/cfg/param/cvi_sdr_bin_SC035HGS
/etc/mtab
/mnt/data/sensor_cfg_OV5647_J1.ini
/mnt/system/ko/cv181x_vo.ko
/etc/init.d/boot
/mnt/system/usr/bin/sample_filterandcsc
/mnt/system/lib/libcipher.so
/mnt/system/lib/libcvi_ive_tpu.so
/mnt/data/sensor_cfg.ini
/mnt/cfg/param/cvi_sdr_bin
/mnt/system/ko/cvi_wiegand.ko
/lib/libatomic.so.1.2.0
/mnt/system/usr/bin/sample_scene_auto
/mnt/system/ko/aic8800_fdrv.ko
/mnt/system/usr/bin/sample_bernsen
/lib/functions/preinit.sh
/mnt/system/usr/bin/sample_and
/bin/ipcalc.sh
/mnt/system/lib/libsys.so
/etc/diag.sh
/mnt/system/lib/libcvi_VoiceEngine.so
/mnt/system/lib/libaacsbrenc2.so
/mnt/system/auto.sh
/usr/bin/duos-pinmux
/mnt/system/ko/cvi_mipi_rx.ko
/mnt/system/lib/libdnvqe.so
/mnt/system/lib/libcvimath.so
/mnt/system/usr/bin/sample_ordstatfilter
